# share/hephaestus/default/overlay/etc/skel/.config/autostart/lupdater.desktop
[Desktop Entry]
Name=lupdater
Description=Start lupdater.py in the background for update notifications.
Exec=lupdater &
Terminal=false
